# root/etc/s6-overlay/s6-rc.d/svc-mod-jellyfin-rclone/run
#!/usr/bin/with-contenv bash

/rclone-linux mount realdebrid: /realdebrid --dir-cache-time 10s --allow-other
