# root/etc/s6-overlay/s6-rc.d/init-mod-rclone/run
#!/usr/bin/with-contenv bash

echo "\
    fuse3" >> /mod-repo-packages-to-install.list

mkdir -p /realdebrid
